# mtprotoproxy.service
[Unit]
Description=High-performance Telegram MTProto proxy
After=network.target

[Service]
Type=simple
# Adjust this path if you cloned the repo elsewhere
WorkingDirectory=/opt/mtprotoproxy
# Use env so it works whether python3 is /usr/bin/python3 or elsewhere
ExecStart=/usr/bin/env python3 mtprotoproxy.py
Restart=always
RestartSec=2

# Resource limits
LimitNOFILE=131072
CPUQuota=80%
Nice=-5
IOSchedulingClass=best-effort
IOSchedulingPriority=4

[Install]
WantedBy=multi-user.target 